# recipes-mapio/mapio-init/mapio-init/usr-local-nvme.mount
[Unit]
Description=NVME disk

[Mount]
What=/dev/nvme0n1
Where=/usr/local/nvme
Type=ext4
Options=defaults

[Install]
WantedBy=multi-user.target
