# recipes-mapio/mapio-init/mapio-init/usr-local-nvme.automount
[Unit]
Description=NVMe automount

[Automount]
Where=/usr/local/nvme

[Install]
WantedBy=multi-user.target
